# systemd/minecraft-server-proxy.socket
[Socket]
ListenStream=0.0.0.0:25565
#ListenDatagram=0.0.0.0:25565
Transparent=true

Accept=no

[Install]
WantedBy=sockets.target
