# systemd/minecraft-server-proxy.socket
[Socket]
ListenStream = 25565
ListenDatagram = 25565

Accept=no
